# Linux/config/systemd/user/default.target.wants/timeupdate.service
[Unit]
Description= run time and date in status bar

[Service]
ExecStart=/usr/bin/timeupdate --no-daemon

[Install]
WantedBy=default.target
